# service_files/app.service
[Unit]
Description=app

[Service]
Type=simple
User=root
ExecStart=/usr/bin/node /var/www/html/server.js
EnvironmentFile=/etc/systemd/system/appEnvVar.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target